# com.bookmark.nav.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bookmark.nav.backend</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/engine/.openclaw/workspace/Projects/BookMark/bookmark-nav-backend/server.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/engine/.openclaw/workspace/Projects/BookMark/bookmark-nav-backend</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PORT</key>
    <string>8080</string>
    <key>NODE_ENV</key>
    <string>production</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>/Users/engine/.openclaw/workspace/Projects/BookMark/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/engine/.openclaw/workspace/Projects/BookMark/logs/launchd.err.log</string>
</dict>
</plist>
